# deploy/systemd/qwen3-tts-predictor.service
[Unit]
Description=Qwen3-TTS Predictor Worker
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
Environment=HOME=%h
Environment=XDG_CONFIG_HOME=%h/.config
Environment=PATH=%h/.cargo/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
Environment=QWEN3_TTS_REPO=kautism/qwen3-tts-rk3588
Environment=QWEN3_TTS_PREDICTOR_BIND=0.0.0.0:9091
Environment=QWEN3_TTS_PREDICTOR_QUANT=q4
EnvironmentFile=-%h/.config/qwen3-tts/systemd.env
WorkingDirectory=%h/.config/qwen3-tts
ExecStartPre=/usr/bin/mkdir -p %h/.config/qwen3-tts
ExecStart=/usr/bin/env qwen3-tts worker -r predictor -b ${QWEN3_TTS_PREDICTOR_BIND} --quant ${QWEN3_TTS_PREDICTOR_QUANT} --repo ${QWEN3_TTS_REPO}
Restart=always
RestartSec=2
NoNewPrivileges=true

[Install]
WantedBy=default.target
